# logrotate/crontab
5 4 * * 0 /usr/sbin/logrotate /home/metaproxy/metaproxy-SUB/logrotate/metaproxy.logrotate
